# default/config/systemd/user/mako-watcher.path
[Path]
PathModified=%h/.config/mako

[Install]
WantedBy=graphical-session.target
